# .config/systemd/wallpaper.service
[Unit]
Description=Set random wallpaper with feh

[Service]
ExecStart=%h/.local/bin/random_wallpaper.sh
